# work/kaga/kagacrawl.timer
[Unit]
Description=Runs kagacrawl

[Timer]
OnCalendar=*-*-* 10:00:00
Unit=kagacrawl.service

[Install]
WantedBy=multi-user.target